# All Scripts and Instructions in Markdown Form
1️⃣ Install MX Linux (clean base)

Flash ISO to USB, reinstall MX Linux (no need for systemd).

After first boot:

sudo apt update && sudo apt upgrade -y
sudo apt install -y git python3-venv python3-pip libcap2-bin ufw

2️⃣ Create virtual env and install OpenCanary

sudo mkdir -p /opt/opencanary
cd /opt/opencanary
sudo python3 -m venv venv
sudo chown -R $USER:$USER venv
source venv/bin/activate
pip install --upgrade pip setuptools==68.0.0
pip install opencanary==0.9.6
deactivate

3️⃣ Initialize config

sudo mkdir -p /etc/opencanaryd
sudo /opt/opencanary/venv/bin/python3 -m opencanaryd --copyconfig
sudo chown -R opencanary:opencanary /etc/opencanaryd

Edit /etc/opencanaryd/opencanary.conf later for modules and syslog target.

4️⃣ Grant Python privileged-port access

sudo setcap 'cap_net_bind_service=+ep' "$(readlink -f /opt/opencanary/venv/bin/python)"

5️⃣ Install SysV service script

sudo nano /etc/init.d/opencanary

#!/bin/sh
### BEGIN INIT INFO
# Provides:          opencanary
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: OpenCanary honeypot (Twisted)
# Description:       Starts OpenCanary using twistd and the opencanary.tac file
### END INIT INFO

NAME="opencanary"
VENV="/opt/opencanary/venv"
TWISTD="$VENV/bin/twistd"
PIDDIR="/var/run/opencanary"
PIDFILE="$PIDDIR/opencanary.pid"
LOGFILE="/var/log/opencanary.log"

# Resolve TAC path from installed package; fallback to source checkout
TAC="$($VENV/bin/python - <<'PY'
import os
try:
    import opencanary
    p = os.path.join(os.path.dirname(opencanary.__file__), "opencanary.tac")
    print(p if os.path.isfile(p) else "")
except Exception:
    print("")
PY
)"
[ -z "$TAC" ] && TAC="/opt/opencanary-src/bin/opencanary.tac"

# Arguments to twistd: use our TAC, log file, and where to write the pid
DAEMON_OPTS="-y $TAC -l $LOGFILE --pidfile $PIDFILE"

do_start() {
    if [ ! -x "$TWISTD" ]; then
        echo "twistd not found at $TWISTD"; return 1
    fi
    if [ ! -f "$TAC" ]; then
        echo "TAC not found at $TAC"; return 1
    fi

    echo "Starting $NAME..."

    # Ensure runtime dirs/log exist on every boot (/var/run is tmpfs)
    mkdir -p "$PIDDIR"
    chown opencanary:opencanary "$PIDDIR"
    touch "$LOGFILE"
    chown opencanary:opencanary "$LOGFILE"

    # Start in background; let twistd create the pidfile (no --make-pidfile here)
    start-stop-daemon --start --oknodo --background \
      --pidfile "$PIDFILE" --chuid opencanary --exec "$TWISTD" -- $DAEMON_OPTS
}

do_stop() {
    echo "Stopping $NAME..."
    if [ -f "$PIDFILE" ]; then
        start-stop-daemon --stop --oknodo --pidfile "$PIDFILE"
        rm -f "$PIDFILE"
    else
        # Fallback if pidfile is missing
        pkill -f "twistd .*opencanary.tac" 2>/dev/null || true
    fi
}

do_status() {
    if [ -f "$PIDFILE" ] && ps -p "$(cat "$PIDFILE" 2>/dev/null)" >/dev/null 2>&1; then
        echo "$NAME running (pid $(cat "$PIDFILE"))"; return 0
    fi
    pgrep -f "twistd .*opencanary.tac" >/dev/null 2>&1 && { echo "$NAME running"; return 0; }
    echo "$NAME not running"; return 3
}

case "$1" in
  start)   do_start ;;
  stop)    do_stop ;;
  restart) do_stop; sleep 1; do_start ;;
  status)  do_status ;;
  *) echo "Usage: /etc/init.d/$NAME {start|stop|restart|status}"; exit 1 ;;
esac
exit $?

Then:

sudo chmod +x /etc/init.d/opencanary
sudo update-rc.d opencanary defaults

6️⃣ Create service user and log dirs

sudo useradd -r -s /usr/sbin/nologin opencanary 2>/dev/null || true
sudo mkdir -p /var/run/opencanary /var/log
sudo touch /var/log/opencanary.log
sudo chown -R opencanary:opencanary /opt/opencanary /etc/opencanaryd /var/run/opencanary /var/log/opencanary.log

7️⃣ Start & verify

sudo service opencanary start
sudo service opencanary status
sudo tail -f /var/log/opencanary.log
sudo ss -tulpen | grep -E ':21|:22|:80' || true

✅ Should show twistd listening on 21, 22, 80.

8️⃣ Firewall & network

sudo ufw allow 21/tcp
sudo ufw allow 22/tcp
sudo ufw allow 80/tcp
sudo ufw enable
sudo ufw status

9️⃣ Auto-start check

sudo reboot
sudo service opencanary status

Should show opencanary running (pid XXXX) automatically.

🔟 Syslog to Kali Purple 

Edit /etc/opencanaryd/opencanary.conf:

"syslog": { "enabled": true, "host": "KALI_IP", "port": 514 }

Then: 

sudo service opencanary restart

11️⃣ Testing your honeypot

From my "victim" Windows device (you will use another device in your Network):

curl http://<HONEYPOT_IP>/
telnet <HONEYPOT_IP> 21
ssh <HONEYPOT_IP> -p 22

Back on MX:

sudo tail -f /var/log/opencanary.log

12️⃣ Log rotation (optional but recommended)

sudo tee /etc/logrotate.d/opencanary >/dev/null <<'EOF'
/var/log/opencanary.log {
  weekly
  rotate 8
  compress
  missingok
  notifempty
  create 0640 opencanary opencanary
  postrotate
    /usr/sbin/service opencanary restart > /dev/null 2>&1 || true
  endscript
}
EOF

